# backend/systemd/hera-node-keep-alive.service
# Configuration file for systemd that keeps the HERA node keep-alive
# daemon running.
#
# Copy this file to /etc/systemd/system/hera-node-keep-alive.service . Then run
# `systemctl enable hera-node-keep-alive` and `systemctl start hera-node-keep-alive`.
#
# This service is meant to be run on hera-node-head.

[Unit]
Description=HERA Node Keep Alive Daemon

[Service]
Type=simple
Restart=always
RestartSec=60
User=hera
Group=hera
ExecStart=/usr/local/bin/hera_node_keep_alive.py

[Install]
WantedBy=multi-user.target
