# monitor-control/systemd/hera-node-cmd-check.service
# Configuration file for systemd that keeps the HERA node command check
# daemon running.
#
# Copy this file to /etc/systemd/system/hera-node-cmd-check.service . Then run
# `systemctl enable hera-node-cmd-check` and `systemctl start hera-node-cmd-check`.
#
# This service is meant to be run on hera-node-head.

[Unit]
Description=HERA Node Command Check Daemon

[Service]
Type=simple
Restart=always
RestartSec=60
User=hera
Group=hera
ExecStart=/usr/local/bin/hera_node_cmd_check.py

[Install]
WantedBy=multi-user.target
